# debian/titan-api.service
[Unit]
Description=Titan-X API Server (FastAPI + uvicorn)
Documentation=https://github.com/titan/titan-x
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/titan-x
Environment=PYTHONPATH=/opt/titan-x/server:/opt/titan-x/core
EnvironmentFile=-/opt/titan-x/.env
ExecStart=/opt/titan/venv/bin/uvicorn server.titan_api:app \
    --host 0.0.0.0 --port 8080 \
    --workers 2 --loop uvloop --http httptools
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

# Resource limits — prevent runaway on shared VPS
MemoryMax=2G
CPUQuota=200%

# Security hardening
ProtectSystem=strict
ReadWritePaths=/opt/titan/data /opt/titan/venv /opt/titan/cuttlefish /tmp
ReadWritePaths=/opt/titan-x/.env /opt/titan-x/docker/ssl
ProtectHome=read-only
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
